# profile.d/image-optim.sh
#!/usr/bin/env bash

export PATH="$PATH:$HOME/vendor/image-optim"
